# systemd/trash-downloads.service
[Unit]
Description=Clean stale Downloads into .trash at boot
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash %h/tools/trash_downloads.sh

[Install]
WantedBy=default.target
