# systemd/user/sleep.target
[Unit]
Description=User level sleep target
StopWhenUnneeded=yes
Wants=/etc/systemd/system/sleep@.service
